# wiki-seer-train/crontab
* * * * * /usr/bin/sh python wikiseertrain/main.py
