# wiki-seer-data/crontab
* * * * * /usr/bin/sh python wikiseerdata/main.py
