# wiki-seer-model/crontab
* * * * * /usr/bin/sh python wikiseermodel/main.py
